# supervisord.conf
[unix_http_server]
file=/var/run/supervisord.sock
chmod=0770

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize -R

[program:crontab]
command=/usr/local/bin/supercronic-linux-amd64 /var/spool/cron/crontabs/www-data


;[supervisorctl]
;serverurl=http://127.0.0.1:9723 ; use a unix:// URL  for a unix socket
[supervisorctl]              ;这个主要是针对supervisorctl的一些配置
serverurl=unix:///var/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[include]
files = /etc/supervisor.d/*.ini